# customize.service
[Unit]
Description=Customize

[Service]
ExecStart=/etc/rc3.d/S01loader
Type=forking
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
